# lib/platforms/freebsd/rc.d/beaverd
#!/bin/sh

# PROVIDE: beaverd
# REQUIRE: LOGIN NETWORKING SERVERS DAEMON FILESYSTEMS
# KEYWORD: shutdown

. /etc/rc.subr

name="beaverd"
rcvar=beaverd_enable
pidfile="/var/run/${name}.pid"
logfile="/var/log/${name}.log"

load_rc_config $name
: ${beaverd_port=1001}
: ${beaverd_home="/usr/local/etc/beaver"}

command_service="beaver -vvvvvv --daemon --home=${beaverd_home} --httpsPort=${beaverd_port}"
command="/usr/sbin/daemon"
command_args="-f -r -H -o ${logfile} -P ${pidfile} -- ${command_service}"

USER=root
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/sbin:/usr/local/bin:/root/bin
PWD=/root
HOME=/root
export NODE_ENV=production

run_rc_command "$1"
